# make-openwrt/openwrt-files/common-files/etc/uci-defaults/99_core
#!/bin/sh

chmod +x /etc/openclash/core/*
exit 0
